# Linux/scripts/services/atf_dga_evaluation.service
[Unit]
Description=running dga_evaluate binary to evaluate the domain as dga or non-dga
After=network.target

[Service]
Type=simple
Restart=always
RestartSec=10
User=attackfence
WorkingDirectory=/opt/attackfence/donaticsInstaller/Linux/scripts/src/
ExecStart=/opt/attackfence/donaticsInstaller/Linux/scripts/src/dga_evaluate

[Install]
WantedBy=multi-user.target
