# Linux/scripts/services/atf_ti_verdict.service
[Unit]
Description=ATF TI VERDICT
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=300
User=attackfence
WorkingDirectory=/opt/attackfence/donaticsInstaller/Linux/scripts/src/
ExecStart=/usr/bin/python3 /opt/attackfence/donaticsInstaller/Linux/scripts/src/analyticsTIReq.py

[Install]
WantedBy=multi-user.target

